# config/mail-helper-celery-beat.service
[Unit]
Description=KI-Mail-Helper Celery Beat Scheduler
After=network.target postgresql.service redis-server.service mail-helper-celery-worker.service
Requires=postgresql.service redis-server.service

[Service]
Type=simple
User=thomas
Group=thomas
WorkingDirectory=/home/thomas/projects/KI-Mail-Helper-Dev
Environment="PATH=/home/thomas/projects/KI-Mail-Helper-Dev/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/thomas/projects/KI-Mail-Helper-Dev/.env.local

# Celery Beat für scheduled Tasks
ExecStart=/home/thomas/projects/KI-Mail-Helper-Dev/venv/bin/celery \
    -A src.celery_app beat \
    --loglevel=info \
    --logfile=/var/log/mail-helper/celery-beat.log \
    --pidfile=/var/run/mail-helper/celery-beat.pid \
    --schedule=/var/run/mail-helper/celerybeat-schedule

# Graceful Shutdown
Restart=always
RestartSec=10s
KillSignal=SIGTERM
TimeoutStopSec=15

# Security Hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource Limits (Beat ist lightweight)
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
